# rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

(
	make clean
	cd /root/Development/ncpa/build
	make build_rpm > /var/www/html/ncpa/build_log.log 2>&1
	cp *.rpm /var/www/html
	make clean
)
